# scripts/launchd/com.mlbpred.scrape_results.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mlbpred.scrape_results</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vinayaks/Desktop/projects/MLB_Predictions/scripts/scrape_results.sh</string>
    </array>

    <!-- 11:00 PM PT — most West Coast games finished by then -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/vinayaks/Desktop/projects/MLB_Predictions/logs/scrape_results.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vinayaks/Desktop/projects/MLB_Predictions/logs/scrape_results_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
